# deploy/pocketbase/cloudflared-pocketbase.service
[Unit]
Description=Cloudflare Tunnel (quick) -> PocketBase 127.0.0.1:8090
After=network-online.target pocketbase.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --url http://127.0.0.1:8090
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
