# deploy/pocketbase/pocketbase.service
[Unit]
Description=PocketBase (Compost Logger backend)
After=network.target

[Service]
Type=simple
User=pocketbase
Group=pocketbase
# Bind to localhost only; Caddy terminates TLS and proxies to this.
# --origins locks CORS to the GitHub Pages app origin. Set here rather than in Caddy
# because PocketBase already emits its own CORS headers; doing it in both places would
# send two Access-Control-Allow-Origin headers and browsers reject that.
ExecStart=/opt/pocketbase/pocketbase serve --http=127.0.0.1:8090 --origins=https://aalkhalifa.github.io
WorkingDirectory=/opt/pocketbase
Restart=always
RestartSec=5
LimitNOFILE=4096

# Basic hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/pocketbase/pb_data

[Install]
WantedBy=multi-user.target
